# goweb.service
[Unit]
Description=goweb

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/usr/local/bin/app

[Install]
WantedBy=multi-user.target
